# cron/tosh.cron
* * * * * [ $(cat /sys/class/power_supply/ACAD/online) -eq 0 ] && /home/mirsella/dev/bin/unpluged
