# private_dot_config/systemd/user/openchamber.service
[Unit]
Description=OpenChamber web server
Documentation=https://github.com/openchamber/openchamber
After=network-online.target graphical-session.target opencode.service
Wants=network-online.target opencode.service
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=%h
EnvironmentFile=%h/.config/opencode/server.env
EnvironmentFile=%h/.config/openchamber/server.env
ExecStart=/usr/bin/openchamber serve --host 127.0.0.1 --port 14097 --foreground
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=graphical-session.target
