# private_dot_config/systemd/user/opencode.service
[Unit]
Description=OpenCode headless server
Documentation=https://opencode.ai/docs
After=network.target

[Service]
Type=simple
WorkingDirectory=%h
EnvironmentFile=%h/.config/telegram.env
ExecStart=opencode serve --port 4096
Restart=on-failure
RestartSec=2s
OOMPolicy=continue
OOMScoreAdjust=-500

[Install]
WantedBy=default.target
